# .config/systemd/user/sync-mail.timer
[Unit]
Description=Sync mail every 2 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min

[Install]
WantedBy=timers.target
